# mainMenu.service
# Save this as /etc/systemd/system/library-menu.service
[Unit]
Description=Library Menu Program
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/pi/library
ExecStart=/usr/bin/python3 /home/pi/library/mainMenu.py
Restart=always
RestartSec=10
# Additional security settings for root service
ProtectSystem=full
PrivateDevices=true

[Install]
WantedBy=multi-user.target

# Commands to set up the service:
sudo chmod 644 /etc/systemd/system/library-menu.service
sudo systemctl daemon-reload
sudo systemctl enable library-menu.service
sudo systemctl start library-menu.service

# To check status:
sudo systemctl status library-menu.service

# To view logs if there are any issues:
# journalctl -u library-menu.service -n 50